# other-scripts/legacy/crontab
@reboot ~/ism-catcher/other-scripts/legacy/rtl_433.sh
