# other-scripts/crontab
@reboot ~/ism-catcher/other-scripts/run.sh init
